# service/luci_swarm/run
#!/bin/sh

setuidgid swarming env HOME=/home/swarming /usr/local/bin/bootstrapswarm -hostname netbsd-arm64-bsiegert
