# service/luci_swarm/run
#!/bin/sh

. /etc/rc.conf

[ -n ${builder_name} ] || exit 1

setuidgid swarming env HOME=/home/swarming /usr/local/bin/bootstrapswarm -hostname ${builder_name}
