# rcd/rc.local
echo -n ' svscanboot'
nohup /usr/pkg/bin/svscanboot &

echo '.'
